# com.jacob.qetrainer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jacob.qetrainer</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>/Users/jp/qe-trainer/bot.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/jp/qe-trainer</string>

    <!-- Start at login and keep it alive if it crashes -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <!-- Don't hammer restart if it keeps failing -->
    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>/Users/jp/qe-trainer/logs/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jp/qe-trainer/logs/launchd.err.log</string>
</dict>
</plist>
